# docs/supervisor-worker.conf.example
; Exemple Supervisor — à copier vers /etc/supervisor/conf.d/ sur Linux
; Adapter le chemin vers artisan et l'utilisateur système (souvent www-data).
; Voir aussi README.md section « Production : Supervisor ».

[program:agrofinance-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/agrofinanceplus/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/agrofinance-worker.log
stopwaitsecs=3600
